# com.ihoka.upwork-search.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ihoka.upwork-search</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ihoka/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/ihoka/ihoka/upwork-search/src/main.ts</string>
    </array>
    <key>StartInterval</key>
    <integer>21600</integer>
    <key>StandardOutPath</key>
    <string>/Users/ihoka/ihoka/upwork-search/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ihoka/ihoka/upwork-search/logs/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ihoka/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CLAUDE_BIN</key>
        <string>/usr/local/bin/claude</string>
        <key>UPWORK_TRIAGE_PROFILE</key>
        <string>/Users/ihoka/ihoka/upwork-search/triage-profile.yaml</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/ihoka/ihoka/upwork-search</string>
</dict>
</plist>
